# infra/leadwa-api.service
[Unit]
Description=Leadwa API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/leadwa/api
Environment="PATH=/var/www/leadwa/api/venv/bin:/usr/bin"
ExecStart=/var/www/leadwa/api/venv/bin/uvicorn main:app --host 127.0.0.1 --port 5002
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
